# system/systemd/autostream_rebooter.path
# This file is part of autostream.
# Copyright (c) 2025 Lo-tech Systems Limited. All rights reserved.

[Unit]
Description=Watch /tmp/rebootrequired and trigger autostream_rebooter.service

[Path]
PathExists=/tmp/rebootrequired
Unit=autostream_rebooter.service

[Install]
WantedBy=multi-user.target
